# registryui.service
[Unit]
Description=CNCF Registry UI
StartLimitInterval=30s
StartLimitBurst=3
Wants=network-online.target
After=network-online.target
Wants=registry.service
After=registry.service

[Service]
User=www
Group=www
Type=simple
WorkingDirectory=/www/server/registry/proxy
Environment=NODE_VERSION=default
Environment=NVM_HOME=/www/server/nvm
Environment=NODE_ENV=production
Environment=PORT=5001
Environment=REGISTRY_ADDR=127.0.0.1
Environment=REGISTRY_PORT=5000
ExecStart=/www/server/nvm/nvm-exec node src/index.js
MemoryMax=128M
Restart=always
RestartSec=5
PrivateTmp=true
ProtectSystem=full
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK AF_UNIX
RestrictNamespaces=true

[Install]
WantedBy=multi-user.target
